# code/104_geräte_dateisysteme_fhs/mnt-disk3.automount
[Unit]
Description=Automount for /mnt/disk3
Wants=mnt-disk3.mount
Before=mnt-disk3.mount

[Automount]
Where=/mnt/disk3

[Install]
WantedBy=multi-user.target
